# templates/default/init.conf.erb
#!/bin/bash
#
# chkconfig: - 85 15
# description: Jakarta Tomcat Java Servlets and JSP server
# processname: tomcat
# pidfile: /var/www/tomcat/tomcat.pid

# Source function library.
. /etc/init.d/functions

TOMCAT_USER=<%= node['tomcat-all']['user'] %>
CATALINA_HOME=<%= node['tomcat-all']['tomcat_home'] %>
CATALINA_PID=$CATALINA_HOME/catalina.pid

case "$1" in
start)
    echo "Starting Tomcat"
    exec su $TOMCAT_USER -c "$CATALINA_HOME/bin/catalina.sh start"
;;

stop)
    echo "Stopping Tomcat"
    exec su $TOMCAT_USER -c "$CATALINA_HOME/bin/catalina.sh stop -force"
;;

restart)
    $0 stop
    sleep 15
    $0 start
;;

status)
    if [ -f "$CATALINA_PID" ]; then
        pid=`cat "$CATALINA_PID"`
        if [ "$pid" != "" ]; then
            if [ -d /proc/$pid ]; then
                # Process is running
                ISRUNNING="1"
                echo "Tomcat is running with pid $pid"
                exit 0
            else
                echo "Tomcat not running (removing stalled PID file)"
                rm -f "$CATALINA_PID"
                exit 1
            fi
        fi
    else
        echo "Tomcat is not running (or $CATALINA_PID is missing)"
        exit 1
    fi
;;

*)
    echo $"Usage: $0 {start|stop|restart|status}"
    exit 1
;;
esac
